# deploy/ansible/roles/basis-agent/templates/basis-agent.service.j2
[Unit]
Description=Basis hypervisor agent
After=network-online.target systemd-modules-load.service
Wants=network-online.target
# VMs are spawned as independent systemd scopes; they survive agent
# restarts, which is the whole point of the design.

[Service]
Type=simple
# `basis_agent_rust_log` (group_vars) selects the level. Default
# enables trace on the dataplane + VM-spawn modules so the
# `journalctl -u basis-agent` bundle includes prefix-route /
# proxy-ARP / GARP / VIP-owner decisions and the cloud-hypervisor
# invocation.
Environment="RUST_LOG={{ basis_agent_rust_log }}"
ExecStart=/usr/local/bin/basis-agent --config /etc/basis/host.yaml
# SIGHUP triggers a config reload. `systemctl reload basis-agent` after
# editing /etc/basis/host.yaml.
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
# Needs root for KVM, VFIO binding, bridge/tap management, and
# systemd-run (to spawn cloud-hypervisor scopes).
User=root

[Install]
WantedBy=multi-user.target
