# deploy/ansible/roles/basis-holod/templates/holod.service.j2
[Unit]
Description=Holo Routing Daemon
Documentation=https://github.com/holo-routing/holo
After=network-online.target
Wants=network-online.target
# Don't restart in a tight loop on systemic failure (config error,
# missing libyang, etc.) — surface it to operators by giving up.
# These keys belong in [Unit], not [Service]; systemd warns and
# ignores them if misplaced.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
# holod must be launched as root: it does its own setuid-drop to
# the `holo` user (configured in holod.toml) AFTER grabbing privileged
# sockets and calling prctl(PR_SET_KEEPCAPS) to preserve caps across
# the drop. Don't set User=/Group= here — it would short-circuit that.
# Don't set NoNewPrivileges=yes either — it makes the keepcaps prctl
# fail with EPERM and panics holod at main.rs:128.
ExecStart=/usr/local/bin/holod --config /etc/holo/holod.toml
Restart=on-failure
RestartSec=5s

# Filesystem hardening. holod only needs its database dir + log dir.
# We keep these (they don't fight the privilege model); the cap and
# NoNewPrivileges restrictions above do, so they're omitted.
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/opt/holo /var/log/holo

[Install]
WantedBy=multi-user.target
