# deploy/ansible/roles/basis-controller/templates/basis-controller.service.j2
[Unit]
Description=Basis hypervisor controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# `basis_controller_rust_log` (group_vars) selects the level.
# Default enables trace on placement / reconcile / VIP-owner
# election so the controller side of the bundle pairs with the
# agent's dataplane trace.
Environment="RUST_LOG={{ basis_controller_rust_log }}"
ExecStart=/usr/local/bin/basis-controller --config /etc/basis/basiscontroller.yaml
Restart=on-failure
RestartSec=5s
# Runs as root so it can read TLS material in /etc/basis/tls (mode 0700).
# The controller binds a network port and writes sqlite under
# /var/lib/basis — no other privileges required. Day-1 tradeoff:
# trivially readable certs over a dynamic-user hardening dance.
User=root

[Install]
WantedBy=multi-user.target
